# tachyon/files/tachyon-master.conf
description "tachyon master"

# Start just after the System-V jobs (rc) to ensure networking and zookeeper
# are started. This is as simple as possible to ensure compatibility with
# Ubuntu, Debian, CentOS, and RHEL distros. See:
# http://upstart.ubuntu.com/cookbook/#standard-idioms

limit nofile 32768 32768

start on stopped rc RUNLEVEL=[2345]
respawn

exec {{ tachyon_home }}/bin/tachyon-wrapper.sh master